# meta-brands/meta-fulan/recipes-bsp/stslave-fw/30-stm-stslave-firmware.rules
ACTION=="add", KERNEL=="st231-0", SUBSYSTEM=="coproc-dev", RUN+="/usr/bin/ustslave /dev/st231-0 /boot/video.elf" 
ACTION=="add", KERNEL=="st231-1", SUBSYSTEM=="coproc-dev", RUN+="/usr/bin/ustslave /dev/st231-1 /boot/audio.elf" 
