# meta-oe/recipes-daemons/vsftpd/vsftpd-3.0.2/vsftpd@.service
[Unit]

[Service]
ExecStart=-/usr/sbin/vsftpd
StandardInput=socket
