# meta-oe/recipes-core/initscripts/initscripts-1.0/nocam.sh
#!/bin/sh
### BEGIN INIT INFO
# Provides:          softcam
# Required-Start:
# Required-Stop:
# Should-Start:      cardserver
# Default-Start:     3 4
# Default-Stop:      0 1 6
# Short-Description: nocam
# Description:       Dummy SoftCAM null-support.
### END INIT INFO

# Placeholder for no cam
exit 0

