# meta-oe/recipes-core/initscripts/initscripts-1.0/nocard.sh
#!/bin/sh
### BEGIN INIT INFO
# Provides:          cardserver
# Required-Start:
# Required-Stop:
# Should-Start:      networking
# Default-Start:     3 4
# Default-Stop:      0 1 6
# Short-Description: nocardserver
# Description:       Dummy card-server null-support.
### END INIT INFO

# Placeholder for no card-server
exit 0

