# meta-oe/recipes-core/busybox/busybox-1.27.2/busybox-telnet.service.in
[Unit]
Description=Telnet Service

[Service]
ExecStart=@sbindir@/telnetd -F
Restart=on-failure

[Install]
WantedBy=multi-user.target
