# meta-oe/recipes-core/busybox/busybox-1.27.2/busybox-inetd.service.in
[Unit]
Description=Internet superserver
Documentation=man:inetd(8)
ConditionPathExists=/etc/inetd.conf

[Service]
ExecStart=@sbindir@/inetd -f
Restart=on-failure

[Install]
WantedBy=multi-user.target

