# meta-oe/recipes-connectivity/connman/connman/connman-env.service.in
[Unit]
Description=Generate options for connection service
Before=connman.service
ConditionKernelCommandLine=root=/dev/nfs
After=syslog.target

[Service]
Type=oneshot
ExecStart=@sbindir@/connmand-env
StandardOutput=null

[Install]
WantedBy=connman.service
