# meta-oe/recipes-core/busybox/busybox-1.22.1/busybox-telnetd@.service
[Unit]

[Service]
ExecStart=-/usr/sbin/telnetd
StandardInput=socket
