# meta-oe/recipes-core/busybox/busybox-1.22.1/busybox-telnetd.socket
[Unit]

[Socket]
ListenStream=23
Accept=yes

[Install]
WantedBy=sockets.target
